# files/etc/service/php-fpm/run
#!/usr/bin/env bash

exec /opt/remi/php${PHP_VERSION}/root/usr/sbin/php-fpm -F
